# supervisord/supervisord.conf
[supervisord]
loglevel=info
nodaemon=true
pidfile=/var/run/supervisord.pid
user=root

[unix_http_server]
file=/tmp/supervisor.sock 

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:confd]
autostart=true
autorestart=true
command=/usr/bin/confd.sh 
logfile=/dev/stdout
loglevel=info
priority=1
user=root

[program:prometheus]
autostart=true
autorestart=true
command=/bin/prometheus -config.file=/etc/prometheus/prometheus.yml
logfile=/dev/stdout
loglevel=info
priority=1000
user=root
